# packaging/rpm/bindizr.service
[Unit]
Description=Bindizr DNS Synchronization Service
After=network.target named.service

[Service]
Type=simple
ExecStart=/usr/bin/bindizr start
WorkingDirectory=/etc/bindizr
Restart=on-failure
RestartPreventExitStatus=1
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
